# site/mysql/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mariadb]
command=/root/run.sh
autorestart=true
user=root