# site/php81/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:cron]
command=/bin/sh /usr/local/data/cron/php81/cron.sh

[program:php-fpm81]
command=php-fpm81 -F