# site/php72/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:cron]
command=/bin/sh /usr/local/data/cron/php72/cron.sh

[program:php-fpm72]
command=php-fpm7 -F